# aquable/rootfs/etc/services.d/aquable/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: AquaBle
# Take down the service when the container is stopped
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
